# dne-paqx-packaging/build/dne-paqx.service
# Copyright © 2017 Dell Inc. or its subsidiaries.  All Rights Reserved
[Unit]
Description=Dell DNE PAQX service
After=network.target docker.service rabbitmq-server.service
[Service]
Type=simple
WorkingDirectory=/opt/dell
Environment=HOSTNAME=%H
ExecStart=/usr/bin/docker-compose -f /opt/dell/cpsd/dne-paqx/image/engineering-standards-service/docker-compose.yml up -d
ExecStart=/usr/bin/docker-compose -f /opt/dell/cpsd/dne-paqx/image/dne-paqx/docker-compose.yml up -d
ExecStop=/usr/bin/docker-compose -f /opt/dell/cpsd/dne-paqx/image/dne-paqx/docker-compose.yml down
ExecStop=/usr/bin/docker-compose -f /opt/dell/cpsd/dne-paqx/image/engineering-standards-service/docker-compose.yml down
[Install]
WantedBy=multi-user.target